# crontab
*   *   *   *   *   /bak.sh push >> /tmp/cronlog
*   *   *   *   *   date >> /tmp/cronlog
